# for build/farpost-bot.conf
[program:farpost-bot]
directory=/home/farpost
command=python3 tgbot.py
autostart=true
autorestart=true
stderr_logfile=/var/log/farpostbot.err.log
stdout_logfile=/var/log/farpostbot.out.log